# dot_void_setup/services/start_stirling-pdf/run
#!/bin/sh

sudo podman start Stirling-PDF
